# private_dot_config/systemd/user/hypr-theme.service
[Unit]
Description=Hypr day/night theme switcher

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/theme-schedule.sh
